# unit/pg_stat_sys.service
[Unit]
Description=pg_stat_sys

[Service]
Type=idle
ExecStart=PYTHON PSC_PATH/pg_stat_sys.py
Restart=always

[Install]
WantedBy=default.target